# deploy/systemd/termcall.service
[Unit]
Description=Termcall incoming-call daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=%h/.local/bin/termcall daemon --incoming open-terminal
Restart=on-failure
RestartSec=5s
NoNewPrivileges=true
PrivateTmp=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6

[Install]
WantedBy=default.target
